# rootskel/etc/network/if-pre-up.d/resolv
#!/bin/sh
#Create resolv.conf - suppress error message in IPE2
/bin/touch /tmp/resolv.conf
